# ansible/roles/galera/templates/etc/systemd/system/mariadb.service.d/limits.conf.j2
[Service]
LimitNOFILE={{ galera_limit_no_file }}
